# packaging/postinst
#!/bin/sh
set -e

# Reload systemd and enable/start system unit if present
if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload || true

  if [ -f /usr/lib/systemd/system/super_v.service ] || [ -f /etc/systemd/system/super_v.service ]; then
    # enable and start the service (tolerate failures)
    systemctl enable --now super_v.service || true
  fi
fi

exit 0
